# service/standbykiller_legacy.service
[Unit]
Wants=mosquitto.service
Description=Standbykiller Legacy -desc -sbpwr W
After=mosquitto.service

[Service]
Type=simple
ExecStart=/bin/bash INSTALLDIR/driver/standbykiller_legacy -dev -desc -sbpwr
Restart=always

[Install]
WantedBy=multi-user.target
